# laravel-php-fpm/provision/supervisor.queue.conf
[program:queue]
command=bash -c "sleep 5 && php /code/application/artisan queue:listen"
autorestart=true
startsecs=60
startretries=10
stdout_events_enabled=true
stderr_events_enabled=true