# laravel-php-fpm/provision/supervisor.queue.conf
[program:queue]
command=php /code/artisan queue:listen
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true